# carbon_relay_ng/supervisord.conf
[supervisord]
nodaemon = true

[program:carbon-relay-ng]
command = /usr/local/bin/carbon-relay-ng
stdout_logfile = /var/log/supervisor/%(program_name)s.log
stderr_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true
